# services/pmscript.service
[Unit]
Description=Gerenciador de Senhas

[Service]
ExecStart=/usr/bin/pmscript.sh

[Install]
WantedBy=multi-user.target
